# deploy/systemd/upscaled-print-agent.service
[Unit]
Description=Upscaled Print Agent
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/upscaled/Upscaled_inv_processing
EnvironmentFile=/etc/upscaled/inventory.env
ExecStart=/usr/bin/node dist/printAgent.js
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
